# examples/secretary.service
[Unit]
Description=Personal AI Secretary - capture server (FastAPI/uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/secretary
EnvironmentFile=/opt/secretary/.env
ExecStart=/opt/secretary/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8095
Restart=always
RestartSec=3
# bind to 127.0.0.1 and expose via your reverse proxy / tunnel (don't open uvicorn to the world)

[Install]
WantedBy=multi-user.target
